# systemd/my-equilibrium.service
[Unit]
Description=Equilibrium Bluetooth Service
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
User=dezi
WorkingDirectory=/home/dezi/dev/MYEquilibrium
ExecStart=/home/dezi/dev/MYEquilibrium/venv/bin/python /home/dezi/dev/MYEquilibrium/main.py
Restart=always
RestartSec=3

# Ensure Bluetooth has time to come up
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
